# WorkerService1/chetch-worker-test.service
[Unit]
Description=chetch-worker-test

[Service]
Type=notify
ExecStart=/usr/sbin/CodeMazeLinuxWorker

[Install]
WantedBy=multi-user.target